# deploy/systemd/forge-training-watchdog.service
[Unit]
Description=LowResource-LLM-Forge Training Watchdog
After=forge-training.service
Wants=forge-training.service

[Service]
Type=simple
WorkingDirectory=%h/projects/LowResource-LLM-Forge
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-%h/.config/forge/training.env
ExecStart=%h/projects/LowResource-LLM-Forge/scripts/training_watchdog.py --service forge-training.service --nan-consecutive-limit 3
Restart=always
RestartSec=10
StandardOutput=append:%h/projects/LowResource-LLM-Forge/artifacts/logs/training_watchdog.log
StandardError=append:%h/projects/LowResource-LLM-Forge/artifacts/logs/training_watchdog.log

[Install]
WantedBy=default.target
