# deploy/systemd/forge-training.service
[Unit]
Description=LowResource-LLM-Forge A100 Training
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/LowResource-LLM-Forge
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-%h/.config/forge/training.env
ExecStart=%h/projects/LowResource-LLM-Forge/scripts/start_a100_training.sh
Restart=on-failure
RestartSec=20
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
